# web/scripts/com.mapledong.china-stock-strategy-daily.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mapledong.china-stock-strategy-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/maple/Projects/quant-traffic-light/web/scripts/daily_update.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/maple/Projects/quant-traffic-light</string>
  <!-- Daily 18:00 Beijing time — set Mac system timezone to Asia/Shanghai -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/maple/Projects/quant-traffic-light/results/daily_update.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/maple/Projects/quant-traffic-light/results/daily_update.log</string>
</dict>
</plist>
